# docker/etc/services.d/pandora/run
#!/bin/execlineb -P

cd /srv/pandora/pandora
fdmove -c 2 1
gunicorn -c gunicorn_config.py wsgi:application
