# docker/etc/services.d/pandora-tasks/run
#!/bin/execlineb -P

cd /srv/pandora/pandora
fdmove -c 2 1
./manage.py celery worker -Q default,celery -n pandora-default --pidfile /run/pandora-tasks.pid --maxtasksperchild 1000 -l INFO
